# image/services.d/docker/run
#!/usr/bin/with-contenv bash

exec dockerd --host=unix:///var/run/docker.sock --log-level=error
