# image/services.d/docker-log-prepare/up
if { mkdir -p /var/log/docker }
if { chown nobody:nogroup /var/log/docker }
chmod 02755 /var/log/docker
